# recipes-thunder-lgi/thunder-4.1.0/wpeframework.service.in
[Unit]
Description=wpeframework
After=lighttpd.service iarmbusd.service pwrmgr.service irmgr.service dsmgr.service sysmgr.service authservice.service slauncher.service @EXTRA_AFTER@
OnFailure=reboot-notifier@%n.service

[Service]
Type=notify
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/var/run/dbus/system_bus_socket
ExecStartPre=-/usr/bin/collect_identification_data.sh
ExecStart=/container/WPEFRAMEWORK/launcher/wpeframework.sh start ${CMD%i}
ExecStop=/container/WPEFRAMEWORK/launcher/wpeframework.sh stop
Restart=no
Delegate=yes

[Install]
WantedBy=multi-user.target
